# deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for WechatToObs
After=network-online.target wechat-obs.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared --no-autoupdate tunnel --config /etc/cloudflared/config.yml run
Restart=on-failure
RestartSec=5
# 隧道凭证与配置目录
ConfigurationDirectory=/etc/cloudflared

[Install]
WantedBy=multi-user.target
