# deploy/wechat-obs.service
[Unit]
Description=WeChat to Obsidian Sync Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/data/wechat-to-obs
Environment=PATH=/data/wechat-to-obs/.venv/bin
ExecStart=/data/wechat-to-obs/.venv/bin/python app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
